# opendkim/rootfs/etc/s6-overlay/s6-rc.d/syslog/run
#!/bin/sh

set -e

chmod a+w /dev/stdout
exec /usr/sbin/rsyslogd -n
